# .local/etc/profile.d/direnv.sh
#!/bin/bash

export DIRENV_WARN_TIMEOUT=10s

eval "$(direnv hook bash)"

