# .local/etc/profile.d/git.sh
GIT_EXTERNAL_DIFF="code --diff --new-window"
