# .local/etc/profile.d/gitsign.sh
export GITSIGN_CONNECTOR_ID=https://github.com/login/oauth
