# .local/etc/profile.d/kubeswitch.sh
complete -F _kubeswitch_autocomplete ks
