# .local/etc/profile.d/kubectl.sh
KUBE_EDITOR="code --wait --new-window"

alias k=kubectl
complete -F __start_kubectl k
