# .local/etc/profile.d/mcfly.sh
#!/bin/bash

eval "$(mcfly init bash)"
export MCFLY_RESULTS=20
export MCFLY_RESULTS_SORT=LAST_RUN
