# .local/etc/profile.d/mount.sh
if ! mount | grep -q " /home/dillen/hidrive "; then
    echo "### Mounting HiDrive"
    mount ~/hidrive/
fi
