# .local/etc/profile.d/mount.sh
if ! mount | grep -q " /home/${USER}/hidrive "; then
    echo "### Mounting HiDrive"
    mount ~/hidrive/
fi
