# .local/etc/profile.d/uniget.sh
export UNIGET_USER=true
alias sudoget="sudo \$(which uniget)"
