# .local/etc/profile.d/uniget.sh
export UNIGET_USER=true
export UNIGET_CACHE=file
alias sudoget="sudo \$(which uniget)"
