# .local/etc/profile.d/whalebrew.sh
#!/bin/bash

export WHALEBREW_INSTALL_PATH=${HOME}/.local/bin
